# server-idle-shutdown.service

[Unit]
Description=Server idle watchdog (shut down VM when server empty too long)
Wants=network-online.target
After=network-online.target docker.service

[Service]
Type=oneshot
EnvironmentFile=/opt/systemd-units/server-scripts/server-idle.env
ExecStart=/opt/systemd-units/server-scripts/server-idle-shutdown.py
User=root
Group=root
Nice=10
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
